# partiels_getter.service
[Unit]
Description=Partiel Getter

[Service]
ExecStart=/path/to/main.py
Restart=on-failure

[Install]
WantedBy=graphical.target