# partiels_getter.service
[Unit]
Description=Partiel Getter

[Service]
ExecStart=bash /path/to/launch.sh
Restart=on-failure

[Install]
WantedBy=graphical.target